# core/deploy/hetzner/systemd/archura-core-backup.service
[Unit]
Description=Back up Archura Core PostgreSQL
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/archura/core/deploy/hetzner
ExecStart=/opt/archura/core/deploy/hetzner/backup-postgres.sh
